# systemd/agent-whatsapp.service
[Unit]
Description=WhatsApp bridge for CLI coding agents
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h
Environment=AGENT_ENV_FILE=__INSTALL_DIR__/.env
ExecStart=__INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/server/gateway.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
